# crontab
* * * * * ~/cafdl.sh
# Mandatory blank line